# systemd/sigmond-storage-trim-timestd.service
# /etc/systemd/system/sigmond-storage-trim-timestd.service
[Unit]
Description=Sigmond storage trim — timestd.events @ 30d
Documentation=https://github.com/mijahauan/sigmond
ConditionPathExists=/var/lib/sigmond/sink.db

[Service]
Type=oneshot
User=root
Group=root

# timestd.events is L2 calibration / metrology data (TOA,
# clock_offset, doppler, propagation_mode, hop count, quality
# grade).  Real science data — no live external consumer exists
# yet, so we retain locally for retrospective analysis.  30 days at
# ~45k rows/day ≈ 1.4M rows / few-hundred MB in pending_uploads —
# well within disk budget on bee1.  Easy to upgrade to a real
# consumer pipeline later without changing this timer; just stop
# the timer if a consumer starts draining the queue itself.
ExecStart=/usr/local/bin/smd storage trim --target-db timestd --max-age 30d --yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sigmond-storage-trim-timestd
