# systemd/sigmond-storage-trim-timestd.service
# /etc/systemd/system/sigmond-storage-trim-timestd.service
[Unit]
Description=Sigmond storage trim — timestd.events @ 7d
Documentation=https://github.com/mijahauan/sigmond
ConditionPathExists=/var/lib/sigmond/sink.db

[Service]
Type=oneshot
User=root
Group=root

# timestd.events is L2 calibration / metrology data (TOA,
# clock_offset, doppler, propagation_mode, hop count, quality
# grade).  No live external consumer drains this queue yet, and
# retention here is low priority — it must not crowd out
# higher-priority data on disk.  7 days at ~45k rows/day ≈ 315k
# rows keeps a useful recent window while bounding sink.db growth
# (the hfdl/wspr trims run at 24h).  If a real consumer pipeline
# starts draining the queue itself, stop this timer.
ExecStart=/usr/local/bin/smd admin storage trim --target-db timestd --max-age 7d --yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sigmond-storage-trim-timestd
